# deployment/supervisord/app.ini
[program:msc102]
command=/pathtoproject/venv/bin/gunicorn -b localhost:8080 -w 4 myapp:app
directory=/pathtoproject
user=youruser
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true