# deployment/supervisord/app.ini
[program:msc102]
command=/opt/msc102/venv/bin/gunicorn -b localhost:8080 -w 4 myapp:app
directory=/opt/msc102
user=youruser
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true